# docker/crontab
* * * * * touch /tmp/test 
* * * * * echo "pulling" && cd /home/jovyan/programming-notebooks && git commit -a -m "user changes" && git pull -s -X ours origin master
