# docker/crontab
* * * * * touch /tmp/test 
* * * * * echo "pulling" && cd /home/jovyan/programming-notebooks && git pull origin master
